# headshot_pipeline/deploy/overseas-vps/flashshot-web.service
# FlashShot — Next.js web app (marketing site + upload/review UI).
#
# `next start` bound to LOOPBACK:3001 behind Caddy (:443 → :3001 for everything
# EXCEPT /api/* and /ws/*, which Caddy routes to flashshot-api). The app MUST be
# built first (`npm run build` in /opt/flashshot/web) — see DEPLOY_OVERSEAS_VPS.md.
#
# Install:  sudo install -m644 flashshot-web.service /etc/systemd/system/
# Enable:   sudo systemctl daemon-reload && sudo systemctl enable --now flashshot-web

[Unit]
Description=FlashShot - Next.js web app
Documentation=file:///opt/flashshot/deploy/overseas-vps/DEPLOY_OVERSEAS_VPS.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=flashshot
Group=flashshot
WorkingDirectory=/opt/flashshot/web
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3001

# `next start` reads the prebuilt .next/ from WorkingDirectory. Uses the pinned
# Node binary + the project-local next CLI (no global next install needed).
#
# Node is installed nvm-ISOLATED under /opt/flashshot/.nvm (Node 20 LTS), NEVER
# to /usr/bin/node — this box is shared (Syntropy/Caddy/Tailscale/fwxt-relay)
# and a global NodeSource install was forbidden as a co-tenant risk. The deploy
# builds a STABLE symlink /opt/flashshot/bin/node → the nvm node, so this unit
# does not hardcode a versioned path. If you reinstall Node, just recreate the
# symlink; this ExecStart stays unchanged.
ExecStart=/opt/flashshot/bin/node node_modules/next/dist/bin/next start -p 3001 -H 127.0.0.1

Restart=on-failure
RestartSec=5

MemoryMax=400M
MemoryHigh=300M

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
