# headshot_pipeline/deploy/overseas-vps/flashshot-api.service
# FlashShot — FastAPI server (Gemini headshot pipeline + Paddle payment).
#
# Runs uvicorn bound to LOOPBACK:8001 behind Caddy (:443 -> :8001 for /api/* and
# /ws/*). It loads env from /opt/flashshot/.env (CWD-relative ".env", read by
# pydantic-settings in server/config.py — that is the single source of truth,
# so this unit does NOT use EnvironmentFile). Generation drives the SiliconFlow
# REST API STATELESSLY - there is no browser
# process to attach to, so this unit has NO dependency on any chrome service.
# (The old flashshot-chrome.service was retired with the Chrome→API pivot.)
#
# Install:  sudo install -m644 flashshot-api.service /etc/systemd/system/
# Enable:   sudo systemctl daemon-reload && sudo systemctl enable --now flashshot-api

[Unit]
Description=FlashShot - FastAPI server (pipeline + Paddle payment)
Documentation=file:///opt/flashshot/deploy/overseas-vps/DEPLOY_OVERSEAS_VPS.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=flashshot
Group=flashshot
WorkingDirectory=/opt/flashshot

# config.py sets env_file=".env", resolved against CWD=/opt/flashshot.
# So /opt/flashshot/.env is the one file to edit for all settings.

# Single worker: the resemblance loop is serialized by the in-process job queue
# (one generation at a time), and 1 worker keeps the shared-box RAM budget
# honest. The SiliconFlow API is stateless per-call, so the worker count is a
# capacity dial, NOT a correctness constraint — raise it only if free RAM on
# preflight.sh allows, and bump MemoryMax/MemoryHigh to match.
ExecStart=/opt/flashshot/.venv/bin/python -m uvicorn server.main:app \
  --host 127.0.0.1 \
  --port 8001 \
  --workers 1

Restart=on-failure
RestartSec=5

# Image processing spikes RAM well past the steady-state ~70 MB: each uploaded
# reference photo is decoded full-res (cv2.imread → multi-MB numpy arrays, plus
# Laplacian/gray/cascade buffers in input_quality.assess_reference_photo), and
# generation loads onnxruntime + the inswapper model + postprocess buffers.
# Observed OOM at anon-rss=712 MB on a 2-photo upload with the old 700 M limit.
# This shared 3.8 GB host has about 2.6 GB available with co-tenants loaded.
# A 1.5 GB ceiling covers the measured image-processing spike while preserving
# at least ~1 GB for Syntropy, Caddy, the web process, and the OS.
MemoryMax=1500M
MemoryHigh=1200M

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# opencv / numpy / onnxruntime / mediapipe / rembg load large .so files; the
# hardening above is as strict as the ML stack tolerates. Anything tighter
# (e.g. RestrictAddressFamilies, SystemCallFilter) tends to break the import.
# Network egress to api.siliconflow.cn is NOT restricted (no IPAddressDeny), which
# is required — generation + judge calls leave the box over HTTPS.
ReadWritePaths=/var/lib/flashshot

[Install]
WantedBy=multi-user.target
