# headshot_pipeline/deploy/com.shanxiang.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shanxiang.tunnel</string>
    <key>ProgramArguments</key>
    <array>
        <string>__CLOUDFLARED__</string>
        <string>tunnel</string>
        <string>--no-autoupdate</string>
        <string>run</string>
        <string>shanxiang</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/shanxiang-tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/shanxiang-tunnel-err.log</string>
</dict>
</plist>
